# packaging/debian_amd64/lib/systemd/system/vault-rest.path
[Unit]
Description=openbank account vault REST service

[Path]
PathChanged=/etc/init/vault.conf

[Install]
WantedBy=vault-rest.service
