# packaging/debian_amd64/lib/systemd/system/vault-unit@.path
[Unit]
Description=openbank account vault tenant %i

[Path]
PathChanged=/etc/init/vault.conf

[Install]
WantedBy=vault.service
