# packaging/debian_amd64/lib/systemd/system/vault.path
[Unit]
Description=openbank account vault services namespace
PartOf=vault.service
ConditionPathExists=/etc/init/vault.conf

[Path]
PathChanged=/etc/init/vault.conf
PathModified=/etc/init/vault.conf
Unit=vault.service

[Install]
WantedBy=vault.service
