# packaging/debian_amd64/lib/systemd/system/vault.path
[Unit]
Description=openbank account vault services namespace
PartOf=vault.service

[Path]
PathChanged=/etc/init/vault.conf

[Install]
WantedBy=vault.service
