# packaging/debian/vault-watcher.path
[Unit]
Description=openbank account vault configuration watcher
Wants=observer.service
After=vault.service

[Path]
PathChanged=/etc/vault/conf.d/init.conf
Unit=vault-watcher.service

[Install]
WantedBy=multi-user.target
